# data/tk.ios-webstack.logrotate-daemonizer/tk.ios-webstack.logrotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
        <key>Label</key>
          <string>logrotate (always at midnight)</string>
        <key>ProgramArguments</key>
          <array>
            <string>/usr/local/sbin/logrotate</string>
            <string>/etc/logrotate.conf</string>
          </array>
        <key>StartCalendarInterval</key>
          <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
          </dict>
        <key>LowPriorityIO</key>
          <true/>
        <key>StandardErrorPath</key>
          <string>/dev/null</string>
        <key>StandardOutPath</key>
          <string>/dev/null</string>
  </dict>
</plist>
